# contrib/coriolis.socket
[Unit]
Description=Coriolis gunicorn socket

[Socket]
ListenStream=/run/coriolis.sock
SocketUser=www-data

[Install]
WantedBy=sockets.target
